# roles/zsh/files/profile.d/fzf.zsh
export FZF_DEFAULT_COMMAND='fd --type f --exclude .git --ex'
export FZF_DEFAULT_OPTS='--color light'
